# crontab
* * * * * bash /home/titan/crons/post/update_ip.sh
